# roleta.service
[Unit]
Description=Roleta Virtual — Servidor Flask/Gunicorn
After=network.target
Wants=network.target

[Service]
Type=simple
User=USER_PLACEHOLDER
WorkingDirectory=APP_DIR_PLACEHOLDER
ExecStart=APP_DIR_PLACEHOLDER/venv/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:5000 \
    --timeout 60 \
    --access-logfile APP_DIR_PLACEHOLDER/instance/access.log \
    --error-logfile  APP_DIR_PLACEHOLDER/instance/error.log \
    wsgi:app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roleta

[Install]
WantedBy=multi-user.target
